# demo/playbooks/hello3.service
[Unit]
Description=Hello Service

[Service]
ExecStart=/opt/bin/hello -p 9470
